# web-panel.service
[Unit]
Description=My Custom Web Panel

[Service]
Type=simple
WorkingDirectory=/media/win/code/web-panel
ExecStart=/home/apkawa/.local/bin/uv run streamlit run app.py --server.port 8501 --server.address 0.0.0.0 --client.toolbarMode hidden
Restart=always

[Install]
WantedBy=default.target
